# dot_config/systemd/user/waybar-watch.service
[Unit]
Description=Watch waybar custom config and auto-merge on change
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=%h/.config/waybar/merge.sh
ExecStart=inotifywait -m -e close_write \
    %h/.config/waybar/config.custom.jsonc \
    %h/.config/waybar/style.custom.css \
    --format %%f \
    | xargs -I{} %h/.config/waybar/merge.sh
Restart=on-failure

[Install]
WantedBy=default.target
